# systemd/gloves-verify.timer
[Unit]
Description=Run gloves verify every 5 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
AccuracySec=30s
Persistent=true
Unit=gloves-verify.service

[Install]
WantedBy=timers.target
